# datastreammanager/vagrant/init.d/clean-topics
### BEGIN INIT INFO
# Provides:         clean-topics
# Required-Start:	$network $local_fs
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO
#! /bin/sh


NAME="clean-topics"
DESC="clean kafka topics"

. /lib/lsb/init-functions

case "$1" in
    start)
    	rm -rf /tmp/kafka-logs/
    	rm -rf /var/lib/zookeeper/version-2/*
	    exit 0
	;;
    stop)
    	exit 0
	;;
    status)
    	exit 0
	;;
  *)
	echo "Usage: kafka {start|stop|status}" >&2
	    exit 1
	;;
esac
